# systemd/dmrlink3-bridge.service
# systemd unit for the DMRlink3 conference-bridge router (bridge.py).
#
# Install:
#   sudo cp systemd/dmrlink3-bridge.service /etc/systemd/system/
#   sudoedit /etc/systemd/system/dmrlink3-bridge.service   # set User/Group and paths
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now dmrlink3-bridge
#   journalctl -u dmrlink3-bridge -f
#
# bridge.py defaults to dmrlink.cfg and bridge_rules in its own directory.
# Pass -c and -b in ExecStart to use different paths. To run the bare protocol
# engine without bridging, point ExecStart at dmrlink.py instead.

[Unit]
Description=DMRlink3 IPSC Conference Bridge Router
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cort
Group=cort
WorkingDirectory=/home/cort/dmrlink3
ExecStart=/home/cort/dmrlink3/venv/bin/python /home/cort/dmrlink3/bridge.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Log to journal (journalctl -u dmrlink3-bridge)
# Use console-timed in [LOGGER] LOG_HANDLERS when running under systemd.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dmrlink3-bridge

[Install]
WantedBy=multi-user.target
