# systemd/dmrlink3-dash.service
# systemd unit for the DMRlink3 web dashboard (dashboard/server.py).
#
# Install:
#   sudo cp systemd/dmrlink3-dash.service /etc/systemd/system/
#   sudoedit /etc/systemd/system/dmrlink3-dash.service   # set User/Group and paths
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now dmrlink3-dash
#   journalctl -u dmrlink3-dash -f
#
# The dashboard reads dashboard/config.py for its settings.
# Copy dashboard/config_sample.py to dashboard/config.py and edit before starting.
# The dashboard reconnects automatically if dmrlink3-bridge restarts; it can be
# enabled/disabled independently.

[Unit]
Description=DMRlink3 Web Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cort
Group=cort
WorkingDirectory=/home/cort/dmrlink3/dashboard
ExecStart=/home/cort/dmrlink3/venv/bin/python /home/cort/dmrlink3/dashboard/server.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Log to journal (journalctl -u dmrlink3-dash)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dmrlink3-dash

[Install]
WantedBy=multi-user.target
